# ip-upstart.conf
# Basic startup script for the ip service
start on (filesystem and net-device-up)

console log

script
	cd /usr/local/ip && gunicorn -u ip -g ip main:app -b 0.0.0.0:5000
end script
